# auto-db-write-job/crontab
*/3 * * * * /usr/local/bin/python /usr/src/script/main.py